# system/paperless/quadlet/paperless-tika.container
[Container]
Image=paperless-tika.image
Pod=paperless.pod

[Service]
Restart=unless-stopped
