# system/paperless/quadlet/paperless-broker.container
[Unit]
PartOf=paperless-webserver.service

[Container]
Image=paperless-broker.image
Pod=paperless.pod

[Service]
Restart=always
