# system/paperless/quadlet/paperless-broker.container
[Container]
Image=paperless-broker.image
Pod=paperless.pod

[Service]
Restart=always
