# system/rclone/systemd/podman.service.d/rclone-socket.conf
[Unit]
Requires=podman-volume-rclone.socket
After=podman-volume-rclone.socket
