# system/hass/quadlet/hass.container
[Unit]
Description=Run Home Assistant
Wants=matter-server.container

[Container]
Image=hass.image
ContainerName=hass
Volume=/var/opt/hass/config:/config
Environment=TZ=Europe/Madrid
Network=host
PodmanArgs=--privileged

[Service]
Restart=always
TimeoutStartSec=900
TimeoutStopSec=70

[Install]
WantedBy=default.target
