# system/rclone/systemd/podman-volume-rclone.socket
[Unit]
Description=Podman Volume Plugin for rclone

[Socket]
ListenStream=/run/podman/plugins/rclone.sock

[Install]
WantedBy=sockets.target
