# system/podman/systemd/podman-image-prune.timer
[Unit]
Description=Prune unused Podman images daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
